# framework/files/usr/lib/systemd/system/elemental-agent-install.service
[Unit]
Description=Elemental Agent Install
Documentation=https://github.com/rancher-sandbox/cluster-api-provider-elemental
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/elemental-agent --debug --install --config /oem/elemental/agent/config.yaml
Restart=on-failure
RestartSec=5
